# demo/systemd/demo-service-start-by-simple-server/demo-service-start.service
[Unit]
Description=Demo Service Start

[Service]
Type=simple
ExecStart=/usr/local/bin/demo-service-start.sh

[Install]
WantedBy=default.target
